# servers_cookbook/rpi_srst_hist/etc/supervisor/conf.d/gw_srst_hist.conf
[program:gw_srst_hist]
command=gw_srst_hist
user=pi
autostart=true
autorestart=false
startsecs=10
startretries=30
