# rpi_mirror/etc/cron.d/usb-mirror-raspbian
20 4	* * *	pi	/usr/local/sbin/usb-mirror-raspbian
