# servers_cookbook/rpi_srst_hist/etc/cron.d/send_to_nas
#
# Regular cron jobs for the backup of this server to Loos NAS server
#

# boot and root fs to Loos NAS, run every 4h at h+10min
10 */4 * * * root /usr/local/sbin/send_fs_to_nas >/dev/null 2>&1
